# bin/observations_passenger
#!/bin/bash 


#	Copy this script to /etc/init.d/
#
#	sudo cp $0 /etc/init.d/
#	sudo update-rc.d observations_passenger defaults
#
#	If updating, will need to run ...
#
#	systemctl daemon-reload


### BEGIN INIT INFO 
# Provides:					Observations passenger
# Default-Start:		 2 3 4 5 
# Default-Stop:			0 1 6 
# Required-Start:		$remote_fs $syslog
# Required-Stop:		 $remote_fs $syslog
# X-Interactive:		 true 
# Short-Description: Start/stop observations web site 
### END INIT INFO 

PATH=/home/jakewendt/.local/bin:/home/jakewendt/.gem/ruby/2.3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
export GEM_HOME=/home/jakewendt/.gem/ruby/2.3
PASSENGER=/home/jakewendt/.gem/ruby/2.3/bin/passenger
ADDRESS=127.0.0.1 
PORT=3000 
RAILS_ENV=production
ENVIRONMENT=production 
APP=/home/jakewendt/observations
#USER=ruby


. /lib/lsb/init-functions 



#	Some of these settings are in the passenger-standalone.json or nginx.config.erg now.

#	KEEP the PORT variable though as it is used in the "status" command.



case "$1" in 
	start) 
		echo "Starting Observations passenger" 
		cd $APP
#		cmd="$PASSENGER start --port $PORT --environment $ENVIRONMENT --daemonize"	# --max-pool-size 2"
#		settings in passenger-standalone.json
		cmd="$PASSENGER start"
		echo $cmd
		$cmd

		#	Modify the iptables AFTER, just in case passenger fails
		#	We end up with MANY rerouting lines in the iptables
		#		iptables --policy INPUT ACCEPT
		#	Be more specific
		#iptables -A INPUT -p tcp --dport 3000 -m conntrack --ctstate NEW,ESTABLISHED -j ACCEPT

		#	Already open, just redirect.
		iptables -t nat -I PREROUTING -p tcp --dport 80 -j REDIRECT --to-ports 3000
		;; 
	stop) 
		echo "Stopping Observations passenger" 
		iptables -t nat -D PREROUTING -p tcp --dport 80 -j REDIRECT --to-ports 3000
		cd $APP
#		cmd="$PASSENGER stop -p $PORT"
#		settings in passenger-standalone.json
		cmd="$PASSENGER stop"
		echo $cmd
		$cmd
		;; 
	status)
		#	From	/lib/lsb/init-functions 
		status_of_proc -p $APP/tmp/pids/passenger.$PORT.pid "" "passenger" && exit 0 || exit $?
		;;
	*) 
		echo "Usage: $0 start|stop|status" >&2 
		exit 3 
		;; 
esac
